# deploy/systemd/elephantbroker-chat-graph-drain.timer
[Unit]
Description=Run ElephantBroker trade-chat graph extraction drain every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=elephantbroker-chat-graph-drain.service

[Install]
WantedBy=timers.target
